# quadlets/searxng/searxng.container
[Unit]
Description=SearxNG
Wants=searxng-redis
After=searxng-redis

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target

[Container]
Image=docker.io/searxng/searxng
ContainerName=searxng
AutoUpdate=registry

Network=searxng.network
HostName=searxng
PublishPort=8080:8080

Volume=searxng-etc:/etc/searxng

Environment=SEARXNG_BASE_URL=https://searxng.example.com/
Environment=UWSGI_WORKERS=4
Environment=UWSGI_THREADS=4
