# pods/mamstack/lazylibrarian/lazylibrarian.container
[Unit]
Description=Lazy Librarian

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target

[Container]
ContainerName=lazylibrarian
Image=lscr.io/linuxserver/lazylibrarian:latest
HostName=lazylibrarian
Pod=MAMstack.pod
AutoUpdate=registry

Volume=lazylibrarian-config:/config
Volume=/volumes/books:/books

Environment=TZ=Etc/UTC
