# quadlets/dashy/dashy.container
[Unit]
Description=Dashboard
After=caddy.service

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target

[Container]
Image=docker.io/lissy93/dashy:$dashy_version
ContainerName=dashy
AutoUpdate=registry

Network=dashy.network
HostName=dashy

Volume=./user-data:/app/user-data

Environment=NODE_ENV=production
